# server/server_pi/systemd/video-streamer.service
[Unit]
Description=Video Streamer for Raspberry Pi (RTP/RTMP)
Documentation=https://github.com/yourorg/video-streamer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rema
Group=rema
WorkingDirectory=/home/rema/server/server_pi

# Environment variables
Environment="PYTHONUNBUFFERED=1"
Environment="PC_IP=192.168.100.41"

# Start the streamer - defaults to RTP mode from config.yaml
# Override mode with --mode rtmp if needed
ExecStart=/usr/bin/python3 -u -m modules.main --config /home/rema/server/server_pi/config.yaml --pc-ip ${PC_IP}

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=video-streamer

# Resource limits (optional)
# Limit memory to 512 MB
# MemoryMax=512M
# MemoryHigh=400M

[Install]
WantedBy=multi-user.target
